# dynamic-dns.cron
*/10 * * * * root /bin/bash /root/dynamic-dns/run.sh
